# build/config/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
priority=1
command=mongod -f /etc/mongodb.conf
stdout_logfile=/var/log/supervisor/mongodb.log
stderr_logfile=/var/log/supervisor/mongodb.log
autorestart=true

#[program:mock-auth-globus-record]
#priority=2
#command=java -jar /usr/local/auth/wiremock-standalone-2.3.1.jar --root-dir=/usr/local/auth/globus --https-port 9999 --proxy-all="https://nexus.api.globusonline.org" --record-mappings --preserve-host-header --verbose
#stdout_logfile=/var/log/supervisor/auth_globus_record.log
#stderr_logfile=/var/log/supervisor/auth_globus_record.log
#autorestart=true

#[program:mock-auth-globus-playback]
#priority=2
#command=java -jar /usr/local/auth/wiremock-standalone-2.3.1.jar --root-dir=/usr/local/auth/globus --https-port 9995
#stdout_logfile=/var/log/supervisor/auth_globus_playback.log
#stderr_logfile=/var/log/supervisor/auth_globus_playback.log
#autorestart=true

[program:shock-server]
priority=3
command=/opt/go/bin/shock-server -conf /opt/run/build/config/shock.conf
stdout_logfile=/var/log/supervisor/shock.log
stderr_logfile=/var/log/supervisor/shock.log
autostart=false
autorestart=true

[program:mysql]
priority=4
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql.log
autorestart=true

[program:init-kbase]
priority=5
command=/opt/run/build/scripts/init_kbase.sh
stdout_logfile=/var/log/supervisor/init_kbase.log
stderr_logfile=/var/log/supervisor/init_kbase.log
autorestart=false
startsecs=0

[program:handle-service]
priority=6
command=/opt/run/build/scripts/start_handle
stdout_logfile=/var/log/supervisor/handle_service.log
stderr_logfile=/var/log/supervisor/handle_service.log
autostart=false
autorestart=false
startsecs=0

[program:handle-manager-service]
priority=6
command=/opt/run/build/scripts/start_handle_manager
stdout_logfile=/var/log/supervisor/handle_manager_service.log
stderr_logfile=/var/log/supervisor/handle_manager_service.log
autostart=false
autorestart=false
startsecs=0

[program:workspace-service]
priority=6
command=/opt/run/build/scripts/start_workspace
stdout_logfile=/var/log/supervisor/workspace_service.log
stderr_logfile=/var/log/supervisor/workspace_service.log
autostart=false
autorestart=false
startsecs=0

[program:nginx]
priority=7
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
autostart=false
autorestart=true
startsecs=0