# build/config/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
priority=1
command=mongod -f /etc/mongodb.conf
stdout_logfile=/var/log/supervisor/mongodb.log
stderr_logfile=/var/log/supervisor/mongodb.log
autorestart=true

[program:shock-server]
priority=2
command=/opt/go/bin/shock-server -conf /opt/run/build/config/shock.conf
stdout_logfile=/var/log/supervisor/shock.log
stderr_logfile=/var/log/supervisor/shock.log
autorestart=true

[program:mysql]
priority=3
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql.log
autorestart=true

[program:init-kbase]
priority=4
command=/opt/run/build/scripts/init_kbase.sh
stdout_logfile=/var/log/supervisor/init_kbase.log
stderr_logfile=/var/log/supervisor/init_kbase.log
autorestart=false
startsecs=0

[program:nginx]
priority=5
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
autorestart=true
startsecs=0